# scripts/cinepi_app.service
[Unit]
Description=CinePi Camera Application
After=multi-user.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/cinepi_app

# Wait for hardware to be ready before starting
ExecStartPre=/bin/sleep 3

# Start the app
ExecStart=/home/pi/cinepi_app/build/cinepi_app

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cinepi_app

# Allow access to hardware devices
SupplementaryGroups=video,input,gpio

# Resource limits for Pi 3A+ (512MB RAM total)
MemoryLimit=256M
CPUQuota=80%

# Security settings
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
