# scripts/cinepi_app.service
[Unit]
Description=CinePi Camera Application
After=network-online.target ssh.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/cinepi_app

# Wait for hardware to be ready before starting
ExecStartPre=/bin/sleep 5

# Start the app
ExecStart=/home/pi/cinepi_app/build/cinepi_app

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cinepi_app

# Allow access to hardware devices
SupplementaryGroups=video,input,gpio

# Resource limits for Pi 3A+ (512MB RAM total)
MemoryLimit=256M
CPUQuota=80%

# Security settings
NoNewPrivileges=true

[Install]
# IMPORTANT: Do NOT enable this in graphical.target
# If SSH cannot start, you lose access to the Pi!
# Instead manually start with: sudo systemctl start cinepi_app
WantedBy=multi-user.target
